# linux/2014/Valentina/crontab
crontab -e
1 * * * * echo $(date)>home
echo $(date)

crontab -e
1 * * * * tar $(date)>home
echo $(date)
