# linux/Valentina/crontab
crontab -e
1 * * * * echo $(date)>home
echo $(date)
